# supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:django]
directory=/app/backend
command=python manage.py runserver 0.0.0.0:8000 --nothreading
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/django.err.log
stdout_logfile=/var/log/supervisor/django.out.log
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nextjs]
directory=/app/frontend
command=/app/frontend/node_modules/.bin/next start -p 3000 -H 0.0.0.0
autostart=true
autorestart=true
startsecs=10
stderr_logfile=/var/log/supervisor/nextjs.err.log
stdout_logfile=/var/log/supervisor/nextjs.out.log
stopasgroup=true
killasgroup=true
environment=NODE_ENV="production"
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
stopasgroup=true
killasgroup=true
